# bsd/custom/rc.local
#! /bin/sh -
# site-specific startup actions, daemons

sysctl -n kern.version | head -1 > /tmp/t1
tail +2 /etc/motd >>/tmp/t1
mv /tmp/t1 /etc/motd
chmod 644 /etc/motd

echo "setting time via ntp" >/dev/console
ntp -s -f chronos >/dev/console 2>&1

echo -n "starting ntpd" >/dev/console
ntpd >/dev/console 2>&1
echo "." >/dev/console
